# deploy/systemd/twinder-deploy.service
[Unit]
Description=Twinder git-poll deploy (one shot)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=%h/.local/bin/twinder-deploy run
